# contestant_containers/files/supervisord-blackarch.conf
[supervisord]
nodaemon=true
#loglevel=trace
serverurl=unix:///run/supervisord.sock
logfile=/var/log/supervisord.log

[inet_http_server]
port=127.0.0.1:9001

[unix_http_server]
file=/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:sshdir]
command=mkdir -p /run/sshd
autostart=true
autorestart=false
exitcodes=0
startsecs=0
priority=1

[program:sshkeys]
command=/usr/bin/ssh-keygen -A
autostart=true
autorestart=false
exitcodes=0
startsecs=0
priority=100
depends_on=sshdir

[program:sshd]
command=/sbin/sshd -D
autostart=true
autorestart=true
restartpause=3
priority=200
depends_on=sshkeys

[program:X11]
command=/sbin/Xvfb :0 -screen 0 1280x800x24
autorestart=true
priority=300

[program:x11vnc]
command=/sbin/x11vnc -localhost -cursor arrow -nopw -display :0 -xkb -noxrecord -noxdamage -wait 5 -shared -forever -nomodtweak
autorestart=true
depends_on=X11
priority=400

[program:xfce]
environment=DISPLAY=":0",HOME="/root",USER="root"
command=/sbin/xfce4-session
autorestart=true
depends_on=X11
priority=500

[program:novnc]
command=/sbin/novnc --vnc localhost:5900 --listen 8080
autorestart=true
depends_on=x11vnc
priority=600
